# files/common/usr/lib/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 01:00:00
