# files/common/usr/lib/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 01:00:00 # BlueBuild Workflow runs at 0600 UTC
                          # Server TZ is Los_Angeles, so this update
                          # runs well after the workflow completes.
